# docker/crontab
# Cron tab for php-redis-admin to collect redis statistics
* * * * * root cd /var/www/html/php-redis-admin/web && /usr/local/bin/php index.php cron/index
